# deploy/systemd/agentgrid-hub.service
[Unit]
Description=AgentGrid Hub
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/agentgrid-hub
ExecStart=/opt/agentgrid-hub/bin/agentgrid-hub --host 127.0.0.1 --port 20181 --db /opt/agentgrid-hub/data/agentgrid-hub.db --web-dir /opt/agentgrid-hub/web
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target

